# rtcback.service
[Unit]
Description=webrtc back service

[Service]
ExecStart=/home/pi/app/pipuup
Restart=always
[Install]
Target=multi-user.target